# .config/systemd/user/check_battery.timer
[Unit]
Description=Check battery

[Timer]
OnCalendar=*-*-* *:0,10:00
Persistent=true

[Install]
WantedBy=timers.target
