# .config/systemd/user/run_rclone_backup.timer
[Unit]
Description=Run rclone mirror to pi

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
