# .config/systemd/user/run_rclone_backup.service
[Unit]
Description=Run rclone mirror

[Service]
Type=oneshot
ExecStart=/home/jamie/scripts/run_rclone_mirror.sh
