# .config/systemd/user/take_screenshot.service
[Unit]
Description=Take screenshot for diary

[Service]
Type=oneshot
ExecStart=/home/jamie/scripts/take_screenshot_diary.sh
